# com.openclaw.grid-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.grid-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/openclaw/trading_env/bin/python3</string>
        <string>-c</string>
        <string>
import sys
import os
os.environ['OKX_API_KEY']='7c1d51b0-0104-476c-90db-c6dff1f0b090'
os.environ['OKX_API_SECRET']='AF8EB679F8AA4CE9A38C5069CB0737A7'
os.environ['OKX_PASSPHRASE']='AbcD@1234'

from okx_client import grid_bot
grid_bot.GridBot().run()
        </string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/openclaw/.openclaw/workspace/trading</string>
    <key>StandardOutPath</key>
    <string>/Users/openclaw/.openclaw/logs/grid_bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/openclaw/.openclaw/logs/grid_bot.error.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
